# .config/systemd/user/grive.timer
[Unit]
Description= Grive Loop(15m)

[Timer]
OnBootSec=5m
OnUnitActiveSec=15m

[Install]
WantedBy=default.target
