# .config/systemd/user/grive.timer
[Unit]
Description= Grive Loop(15m)

[Timer]
OnActiveSec=15m

[Install]
WantedBy=default.target
