# .config/systemd/user/network-online.target
[Unit]
Description=Network Online Target
