# .config/systemd/user/touchegg.service
[Unit]
Description=Touchegg
Requires=wm.target

[Service]
ExecStart=/usr/bin/touchegg
Restart=on-abort

[Install]
WantedBy=wm.target
