# .config/systemd/user/wall_looper.service
[Unit]
Description=Wallpaper Looper
Requires=wm.target

[Service]
ExecStart=/bin/bash -c "%h/bin/systemd-bin/wall_looper.sh %h/Media/Pictures/i3/bg/ 600"

[Install]
WantedBy=wm.target
