# .config/systemd/user/wm.target
[Unit]
Description=Window Manager
After=default.target
